# cron.d/100syncarchive_volatile
#!/bin/bash

. `dirname $0`/../config.sh

test -d ${archivedir} || mkdir -p ${archivedir}
cd ${archivedir}

# volatile archive
#
for dist in stable
do
  echo "retrieve Release(.gpg) $dist-volatile"
  test -d volatile/${dist}-volatile || mkdir -p volatile/${dist}-volatile
  (cd volatile/${dist}-volatile &&
      wget -q -N ${volatile_ftpsite}/dists/${dist}/volatile/Release &&
      wget -q -N ${volatile_ftpsite}/dists/${dist}/volatile/Release.gpg )

    foo=\$arch_${dist//-/_}
    for arch in `eval echo $foo`
    do
      for part in ${parts}
      do
	echo retrieve Packages $dist-volatile/$part/$arch
	test -d volatile/${dist}-volatile/${part}/binary-${arch} || mkdir -p volatile/${dist}-volatile/${part}/binary-${arch}
	(cd  volatile/${dist}-volatile/${part}/binary-${arch} &&
	    wget -q -N ${volatile_ftpsite}/dists/${dist}/volatile/${part}/binary-${arch}/Packages.gz)
      done
      echo retrieve Contents $dist-volatile/$arch
      (cd volatile/${dist}-volatile &&
	  wget -q -N ${volatile_ftpsite}/dists/${dist}/volatile/Contents-${arch}.gz)

    done
    for part in ${parts}
    do
        echo retrieve Sources $dist-volatile/$part
        test -d volatile/${dist}-volatile/${part}/source || mkdir -p volatile/${dist}-volatile/${part}/source
        (cd volatile/${dist}-volatile/$part/source &&
            wget -q -N ${volatile_ftpsite}/dists/${dist}/volatile/${part}/source/Sources.gz)
    done
done
