# cron.d/200process_archive
#!/bin/bash

. `dirname $0`/../config.sh

cd "$topdir"

date
./bin/parse-packages
date
./bin/parse-sources
date
./bin/parse-contents
date
