# cron.d/400changelogs_map
#!/bin/sh

. `dirname $0`/../config.sh

cd "$topdir"

date
./bin/create_changelogs_map "$filesdir/db/" "$filesdir/changelogs-map"
date
