# cron.d/600index_pages
#!/bin/bash

. `dirname $0`/../config.sh

cd "$topdir"

date
./bin/create_index_pages
date
