# cron.d/600prepare_newpkg
#! /bin/bash

. `dirname $0`/../config.sh

packagesdir=${filesdir}/packages
test -d ${packagesdir} || mkdir -p ${packagesdir}

LANG=C cat ${filesdir}/db/package_names_unstable.txt \
    | sort | uniq > ${packagesdir}/package_names_unstable.$(date -I)

cd $topdir

./bin/newpkg_info unstable
