# cron.d/700install_static
#!/bin/sh

set -e

. `dirname $0`/../config.sh

test -d "$htmldir" || mkdir -p "$htmldir"

cd "$topdir" && ttree --perl5lib=${libdir} --load_perl --eval_perl \
                   -f ./conf/ttreerc
