# cron.d/700install_static
#!/bin/sh

set -e

. `dirname $0`/../config.sh

test -d "$htmldir" || mkdir -p "$htmldir"

cd "$topdir" && ttree -f ./conf/ttreerc
